# systemd/efinder-firstboot.service
[Unit]
Description=eFinder first-boot configuration
# Need NetworkManager up so we can `nmcli con add` AP and USB profiles.
# Don't wait for network-online -- on first boot we *create* the network,
# so network-online.target may never fire.
After=NetworkManager.service
Wants=NetworkManager.service
# Once we mark ourselves done, this unit becomes a no-op (but the
# ConditionPathExists makes the no-op cheap rather than a script invocation).
ConditionPathExists=!/var/lib/efinder/firstboot.done

[Service]
Type=oneshot
ExecStart=/opt/efinder/scripts/firstboot.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
# First boot can take 30-60 seconds (waiting for NM, generating creds,
# detecting hardware). Default timeout is 90s which should be plenty,
# but bump it to be safe.
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
