# systemd/efinder-webui.service
[Unit]
Description=eFinder web UI
After=efinder.service network.target
Wants=efinder.service

[Service]
Type=simple
User=efinder
Group=efinder
WorkingDirectory=/opt/efinder

Environment=PYTHONUNBUFFERED=1
Environment=EFINDER_CONFIG=/etc/efinder/efinder.conf
Environment=FLASK_APP=webui.app

# Bind to port 80 on all interfaces. Needs CAP_NET_BIND_SERVICE since
# we run as the unprivileged efinder user.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

ExecStart=/opt/efinder/venv/bin/python -m webui.app

Restart=on-failure
RestartSec=5

# The web UI runs the efinder-update script via sudo. We don't grant
# the efinder user blanket sudo; instead a sudoers drop-in scoped to
# just /usr/local/bin/efinder-update is installed by install.sh.

NoNewPrivileges=false   # need sudo to work (drops to it via sudoers rule)
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
